# scripts/.config/systemd/user/graphify-sync.service
[Unit]
Description=Rebuild graphify-memory + graphify-personal graphs and push ~/raw

[Service]
Type=oneshot
ExecStart=%h/dotfiles/scripts/graphify-sync.sh
Nice=10
TimeoutStartSec=15min
