# scripts/.config/systemd/user/memory-stale-report.service
[Unit]
Description=Generate Claude Code memory stale-candidates report (weekly)

[Service]
Type=oneshot
ExecStart=%h/dotfiles/scripts/memory-stale-report.sh
Nice=15
TimeoutStartSec=10min
